# scripts/launch_agents/com.ordertech.localdev.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.ordertech.localdev</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>/Users/mosawi/DATATECH/OrderTech/scripts/local_up_only_server.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/mosawi/DATATECH/OrderTech</string>

  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/mosawi/DATATECH/OrderTech/.logs/local_server.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/mosawi/DATATECH/OrderTech/.logs/local_server.log</string>
</dict>
</plist>
